# target/linux/qualcommax/ipq807x/base-files/etc/hotplug.d/firmware/10-ath10k-caldata
#!/bin/sh

[ -e /lib/firmware/$FIRMWARE ] && exit 0

. /lib/functions/caldata.sh

board=$(board_name)

case "$FIRMWARE" in
"ath10k/pre-cal-pci-0001:01:00.0.bin")
	case $board in
	linksys,mx5300)
		caldata_extract "0:art" 0x33000 0x2f20
		ath10k_remove_regdomain
		;;
	esac
	;;
*)
	exit 1
	;;
esac
